# xml_rpc_server/xml_rpc_server.conf
#/etc/init/xml_rpc_server.conf
description "xmk_rpc_server"
start on stopped rc RUNLEVEL=[2345]
respawn
exec python3 /mnt/hgfs/Development/Odoo/dmpi_fresh/xml_rpc_server/xml_rpc_server.py